# service/alertmanager/run
#!/bin/sh

exec chpst -U nobody:nogroup /bin/alertmanager --config.file /etc/prometheus/alertmanager.conf
