# service/mesos-exporter/run
#!/bin/sh

exec chpst -U nobody:nogroup /bin/mesos-exporter
